# mybot4qq.service
[Unit]
Description=MyBot4QQ Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/heke1228/application/MyBot4QQ

ExecStart=/usr/bin/uv run main.py

Restart=on-failure
RestartSec=10
StartLimitBurst=3

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mybot4qq

[Install]
WantedBy=multi-user.target
